# templates/config/systemd/user/i3-session.target
[Unit]
Description=i3 is up and running
BindsTo=graphic-session.target
